# files/firstboot/ublue-firstboot-date.service
[Unit]
Description=Record Dakota first-boot installation date
ConditionPathExists=!/var/lib/ublue-os/install-date
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'mkdir -p /var/lib/ublue-os && date +%%Y-%%m-%%d > /var/lib/ublue-os/install-date'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
