# files/service-overrides/flatpak-preinstall.service.d/condition-check.conf
[Service]
# The flatpak preinstall subcommand was added after 1.16.x.
# Guard with ExecCondition so the service exits cleanly (not failed)
# on images that ship an older flatpak, rather than logging errors on
# every boot.  Drop this override once the minimum flatpak version
# in fdsdk is >= the version that introduced 'flatpak preinstall'.
ExecCondition=/bin/sh -c 'flatpak help 2>&1 | grep -q preinstall'
